# scripts/notifications/com.smartstate.recruiter-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smartstate.recruiter-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/isaacmarks/Desktop/smartstate-recruitment-pipeline-clickup-set-up/scripts/sync/sync_linkedin_recruiter.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/isaacmarks/Desktop/smartstate-recruitment-pipeline-clickup-set-up</string>

    <!-- Run daily before the 9:00 AM cloud nonresponder job -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/smartstate_recruiter_sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/smartstate_recruiter_sync.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
