# scripts/notifications/com.smartstate.replies.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smartstate.replies</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/isaacmarks/Desktop/smartstate-recruitment-pipeline-clickup-set-up/scripts/notifications/run_all.py</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>WorkingDirectory</key>
    <string>/Users/isaacmarks/Desktop/smartstate-recruitment-pipeline-clickup-set-up/scripts/notifications</string>

    <key>StandardOutPath</key>
    <string>/tmp/smartstate_replies.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/smartstate_replies_error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
